# build/cl-systemd/cl-cert.timer
#
# Copyright 2019 Brightgate Inc.
#
# This Source Code Form is subject to the terms of the Mozilla Public
# License, v. 2.0. If a copy of the MPL was not distributed with this
# file, You can obtain one at https://mozilla.org/MPL/2.0/.
#


[Unit]
Description=Run cl-cert periodically

[Timer]
OnCalendar=*-*-* 1/6:00:00
Persistent=true

[Install]
WantedBy=timers.target

