# build/cl-systemd/cl-obs.timer
#
# Copyright 2020 Brightgate Inc.
#
# This Source Code Form is subject to the terms of the Mozilla Public
# License, v. 2.0. If a copy of the MPL was not distributed with this
# file, You can obtain one at https://mozilla.org/MPL/2.0/.
#


[Unit]
Description=Run cron-cl-obs periodically

[Timer]
OnCalendar=*-*-* 1,5,9,13,17,21:00:00
Persistent=true

[Install]
WantedBy=timers.target

