# build/debian-deb/brightgate-appliance.service
#
# Copyright 2018 Brightgate Inc.
#
# This Source Code Form is subject to the terms of the Mozilla Public
# License, v. 2.0. If a copy of the MPL was not distributed with this
# file, You can obtain one at https://mozilla.org/MPL/2.0/.
#


[Unit]
Description=Brightgate appliance starter
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/com.brightgate/bin/ap-ctl start all

[Unit]
Requires=ap.mcp.service
After=ap.mcp.service

[Install]
WantedBy=multi-user.target

